# docker/supervisord.conf
; =============================================================================
; Supervisor配置文件 - 管理Linux微信+VNC+支付回调服务
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; =============================================================================
; DBus系统总线（微信桌面通知需要）
; =============================================================================
[program:dbus]
command=/bin/bash -c 'mkdir -p /run/dbus && dbus-daemon --system --nofork'
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus-error.log

; =============================================================================
; 微信支付回调服务
; =============================================================================
[program:wechat-pay-callback]
command=python3 /app/main.py -c /app/config.yaml
autostart=true
autorestart=true
priority=20
environment=DISPLAY=":1",DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dbus/system_bus_socket",PYTHONUNBUFFERED="1",CALLBACK_URL="%(ENV_CALLBACK_URL)s",CALLBACK_SECRET_KEY="%(ENV_CALLBACK_SECRET_KEY)s",DECRYPT_PROJECT_DIR="%(ENV_DECRYPT_PROJECT_DIR)s",WECHAT_DB_DIR="%(ENV_WECHAT_DB_DIR)s",SESSION_DB_PATH="%(ENV_SESSION_DB_PATH)s"
stdout_logfile=/var/log/supervisor/wechat-pay.log
stderr_logfile=/var/log/supervisor/wechat-pay-error.log
redirect_stderr=false

; =============================================================================
; 微信消息网关（可选，根据需要启用）
; =============================================================================
; [program:message-gateway]
; command=python3 /app/message_gateway.py
; autostart=false
; autorestart=true
; priority=25
; environment=DISPLAY=:1,PYTHONUNBUFFERED=1
; stdout_logfile=/var/log/supervisor/message-gateway.log
; stderr_logfile=/var/log/supervisor/message-gateway-error.log

; =============================================================================
; 微信消息网关（可选，根据需要启用）
; =============================================================================
; [program:message-gateway]
; command=python3 /app/message_gateway.py
; autostart=false
; autorestart=true
; priority=25
; environment=DISPLAY=:1,PYTHONUNBUFFERED=1
; stdout_logfile=/var/log/supervisor/message-gateway.log
; stderr_logfile=/var/log/supervisor/message-gateway-error.log

; =============================================================================
; 事件监听器
; =============================================================================
[eventlistener:process_exit]
command=/bin/bash -c 'echo "Process exited: $1" >> /var/log/supervisor/event.log'
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
